# debian/prerm
#!/bin/sh
/etc/init.d/chilli stop

if [ "$1" = purge ]; then
     		update-rc.d chilli remove
fi

